# crontab/admin
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

#0 8 * * * root /usr/bin/python3 -u /home/ubuntu/nykaa_scripts/feed_pipeline/popularity.py --fetch-omniture --preprocess --popularity --startdate -1 -n0 -y 2>&1 >> /tmp/popularity.log
* * * * * root /usr/bin/python /home/ubuntu/nykaa_scripts/feed_pipeline/indexSchedulePricesOptimized.py --threads 4 --batch_size 100 2>&1 >> /var/log/nykaa/schedule_price_cron.log

# 30 1 * * * root python /home/ubuntu/nykaa_scripts/feed_pipeline/catalogPipelineMultithreaded.py --importattrs -u "http://www.nykaa.com/media/feed/master_feed_gludo.csv"  >> /var/log/nykaa/catalog_pipeline.log

0 0 * * * root /usr/bin/rsync -avzhe ssh /nykaa/adminftp/ ubuntu@52.77.36.71:/nykaa/adminftp
0 0 * * * root /usr/bin/rsync -avzhe ssh /nykaa/adminftp/ ubuntu@13.228.118.41:/nykaa/adminftp
0 0 * * * root /usr/bin/rsync -avzhe ssh /nykaa/adminftp/ ubuntu@qa-api.nyk00-int.network:/nykaa/adminftp

#0 0 * * * root /usr/bin/python /nykaa/scripts/feed_pipeline/cron_runner.py
0 23 * * * root /usr/bin/python /nykaa/scripts/findMismatches.php
0 0 * * * root /usr/bin/python /nykaa/scripts/feed_pipeline/health_check.py -m  
0 1 * * * root /usr/local/bin/aws s3 sync /nykaa/adminftp s3://prod-nykaafeeds/search_autocomplete/ 
0 * * * * root /usr/local/bin/aws s3 cp /data/nykaa/master_feed_gludo.csv  s3://nykaadp-feeds/master_feed_gludo/master_feed_gludo.csv
0 23 * * * root /usr/bin/python /nykaa/scripts/feed_pipeline/product_history_datapipeline.py
