# crontab/dev
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

45 */4 * * * root /home/ubuntu/nykaa_scripts/findMismatches.php -d -e >> /var/log/nykaa/mismatch.log
#0 * * * * root cd /home/ubuntu/nykaa_scripts/feed_pipeline; python indexScheduledPrices.py >> /var/log/nykaa/schedule_price_cron.log

0 0 * * * root /usr/bin/python /nykaa/scripts/feed_pipeline/cron_runner.py 
0 0 * * * root /usr/bin/python /nykaa/scripts/feed_pipeline/health_check.py -m  

